# super/proxys.ini
[program:proxys]
directory=/home/IPProxyTool-master/
command=python runspider.py
autostart=true
user=root
startsecs=10
autorestart=true
startretries=3
priority=999
redirect_stderr=true
stdout_logfile=/home/proxys.log
loglevel=error
