# super/ss.ini
[program:ss]
directory=/root/kcptun/shadowsocks/
command=python /root/kcptun/shadowsocks/server.py
autostart=true
startsecs=10
autorestart=true
startretries=10
priority=999
redirect_stderr=true
stdout_logfile=/home/ss.log
loglevel=error
