# super/web.ini
[program:web]
directory=/home/wwwroot/192.3.244.150/svnwebso/
command=/usr/local/bin/gunicorn -w 5 webso.wsgi:application -b 0.0.0.0:8021 --reload
autostart=true
user=root
startsecs=10
autorestart=true
startretries=3
priority=999
redirect_stderr=true
stdout_logfile=/home/web.log
loglevel=error
